# docker/root/etc/s6-overlay/s6-rc.d/svc-app/run
#!/usr/bin/with-contenv bash

echo -e "\nApp is starting!"
export NODE_ENV=production
cd /app || exit

# Run build first
exec \
	  s6-setuidgid abc /usr/bin/node $NODE_ARGS /app/main.js